# device/temp-sensor.service
[Unit]
Description=Temperature and Humidity Sensor (Prometheus Exporter)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
Environment="SENSOR_LOCATION=default"
Environment="PYTHONUNBUFFERED=1"
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 /home/pi/temp-sensor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
